# at2-ns/at2-ns.service
[Unit]
After=network.target

[Service]
ExecStart=%h/bin/at2-ns
Restart=on-failure

[Install]
WantedBy=multi-user.target
