# dns/at2-dns.service
[Unit]
After=network.target

[Service]
ExecStart=%h/bin/dns
Restart=on-failure

[Install]
WantedBy=multi-user.target
